# agent/spot-agent.service
[Unit]
Description=AWS Spot Optimizer Agent
After=network.target
Documentation=https://github.com/your-repo/spot-optimizer

[Service]
Type=simple
User=root
WorkingDirectory=/opt/spot-optimizer/agent

# Environment variables - UPDATE THESE!
Environment="CENTRAL_SERVER_URL=http://your-central-server:5000"
Environment="CLIENT_TOKEN=your-client-token-here"
Environment="HEARTBEAT_INTERVAL=30"

# Optional: Set if agent already registered
# Environment="AGENT_ID=your-agent-id"

ExecStart=/usr/bin/python3 /opt/spot-optimizer/agent/spot_agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
